# packer/salt/resources/supervisor-mailcatcher.conf
[program:mailcatcher]
command=/usr/local/bin/mailcatcher