# packer/salt/solr5/supervisor.conf
[program:solr5]
command = /opt/solr-5.5.0/bin/solr start -f
user = nobody