# packer/salt/solr5/supervisor.conf
[program:solr5]
command = /opt/solr-{{ solr5_settings.version }}/bin/solr start -f
user = nobody