# packer/salt/virtuoso-opensource-7/supervisor.conf
[program:virtuoso-opensource-7]
command = /usr/bin/virtuoso-t +wait +configfile /etc/virtuoso-opensource-7/virtuoso.ini