# packer/salt/mailcatcher/supervisor.conf
[program:mailcatcher]
command=/usr/local/bin/mailcatcher -f --http-ip 0.0.0.0
user=nobody